# ladon/systemd/node-exporter-collectors.timer
[Unit]
Description=Run Node Exporter Custom Collectors Periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
