# olympus/ladon/systemd/ladon.service
[Unit]
Description=Ladon Monitoring Stack
Requires=docker.service
After=docker.service network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/todd/ML-Lab/Olympus/olympus/ladon
Environment=COMPOSE_PROJECT_NAME=ladon

# Create necessary directories and set permissions
ExecStartPre=/bin/mkdir -p /var/lib/node_exporter/textfile_collector
ExecStartPre=/bin/chown root:root /var/lib/node_exporter/textfile_collector

# Start collectors
ExecStartPre=/bin/bash -c '/home/todd/ML-Lab/Olympus/olympus/ladon/node-exporter-textfile/mdstat.sh'
ExecStartPre=/bin/bash -c '/home/todd/ML-Lab/Olympus/olympus/ladon/node-exporter-textfile/lvm.sh'

# Start the monitoring stack
ExecStart=/usr/bin/docker compose up -d

# Proper shutdown
ExecStop=/usr/bin/docker compose down

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
